# config/.misc/scripts/openwrt/old/tl-wr1043nd/files/etc/inittab
::sysinit:/etc/init
ttyS0::once:/bin/login -f root
