# config/config/etc/inittab
::sysinit:/etc/init
::restart:/usr/sbin/poweroff
::shutdown:/usr/sbin/poweroff
::ctrlaltdel:/usr/sbin/poweroff
tty1::respawn:/usr/bin/login -f root
tty2::respawn:/usr/bin/login -f user
tty3::respawn:/usr/bin/login -f nobody
tty4::respawn:/usr/bin/login -f proxy
tty5::respawn:/usr/bin/login -f jail
