# config/config/etc/inittab
::sysinit:/etc/init
::restart:/sbin/poweroff
::shutdown:/sbin/poweroff
::ctrlaltdel:/sbin/poweroff
ttyS0::respawn:/bin/login -f root
tty1::respawn:/bin/login -f root
tty2::respawn:/bin/login -f user
tty3::respawn:/bin/login -f nobody
tty4::respawn:/bin/login -f proxy
tty5::respawn:/bin/login -f jail
