# config/.misc/scripts/openwrt/old/tl-mr3020/files/etc/inittab
::sysinit:/etc/init
ttyATH0::once:/bin/login -f root
